# systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper from unsplash

[Service]
ExecStart=%h/wallpaper-unsplash.sh

[Install]
WantedBy=default.target
